# crontab
@reboot /home/user/startup.sh
20 10 * * * /home/user/restart-broadcast.sh
